# knockd.service
[Unit]
Description=Port-Knocking Daemon
After=network.target

[Service]
ExecStart=/usr/bin/knockd

[Install]
WantedBy=multi-user.target
